# systemd/llmnews.service
[Unit]
Description=LLM News Fetcher
After=network.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/root/llmnews
Environment=PYTHONUNBUFFERED=1
ExecStart=/root/llmnews/dev.sh run
StandardOutput=append:/var/log/llmnews.log
StandardError=append:/var/log/llmnews.log

[Install]
WantedBy=multi-user.target