# systemd/llmnews.timer
[Unit]
Description=Run LLM News Fetcher daily at 1am UTC
Unit=llmnews.service

[Timer]
OnCalendar=*-*-* 01:00:00 UTC
Persistent=true

[Install]
WantedBy=timers.target